# deploy/polla.service
[Unit]
Description=Polla — World Cup live standings
After=network.target mysql.service

[Service]
Type=simple
# Adjust user and paths to your server:
User=www-data
WorkingDirectory=/var/www/html/polla
EnvironmentFile=/var/www/html/polla/.env
ExecStart=/var/www/html/polla/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
