# deploy/polla.service
[Unit]
Description=Polla — World Cup live standings
After=network.target mysql.service

[Service]
Type=simple
# Adjust user and paths to your server:
User=polla
WorkingDirectory=/opt/polla
EnvironmentFile=/opt/polla/.env
ExecStart=/opt/polla/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
